# system/service/adsb-feeders.target
[Unit]
Description=ADS-B Feeder Services

[Install]
WantedBy=multi-user.target
